# deploy/qwen-image-api/qwen-image-api.service
[Unit]
Description=Qwen-Image-2512 API (FP8 via TorchAO)
After=network-online.target nvidia-persistenced.service
Wants=network-online.target

[Service]
Type=simple
User=pablo
Group=pablo
WorkingDirectory=/home/pablo/Projets/qwen3.5-122/deploy/qwen-image-api

EnvironmentFile=/home/pablo/Projets/qwen3.5-122/deploy/qwen-image-api/.env

ExecStart=/home/pablo/Projets/qwen3.5-122/deploy/qwen-image-api/.venv/bin/python -m uvicorn app:app \
    --app-dir /home/pablo/Projets/qwen3.5-122/deploy/qwen-image-api \
    --host 0.0.0.0 \
    --port 18291

ExecStartPre=/bin/mkdir -p /home/pablo/qwen-image-outputs /home/pablo/Projets/qwen3.5-122/deploy/qwen-image-api/logs

Restart=on-failure
RestartSec=10

# Memory/resource limits
LimitNOFILE=65536
OOMScoreAdjust=-500

# Logging to journald (also accessible via journalctl -u qwen-image-api)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qwen-image-api

[Install]
WantedBy=multi-user.target
